# ecoflow-ac-notify.timer
[Unit]
Description=Run EcoFlow AC notify every minute

[Timer]
OnBootSec=60s
OnUnitActiveSec=300s
Persistent=true

[Install]
WantedBy=timers.target
